# systemd/logging-dashboard.service
[Unit]
Description=Logging Dashboard (Memory Optimized)
Documentation=https://github.com/cleptogk/logger
After=network.target enhanced-logging-api.service
Wants=enhanced-logging-api.service

[Service]
Type=simple
User=logserver
Group=logserver
WorkingDirectory=/opt/logging
Environment=PYTHONPATH=/opt/logging/logger
EnvironmentFile=-/opt/logging/config/environment

# Memory-optimized Gunicorn command for dashboard
ExecStart=/opt/logging/venv/bin/gunicorn \
    --config gunicorn-dashboard.conf.py \
    --bind 0.0.0.0:8081 \
    --workers 1 \
    --worker-class eventlet \
    --worker-connections 50 \
    --max-requests 200 \
    --max-requests-jitter 20 \
    --timeout 60 \
    --keepalive 2 \
    --preload \
    --access-logfile - \
    --error-logfile - \
    --log-level info \
    --proc-name logging-dashboard \
    dashboard.dashboard_app:app

ExecReload=/bin/kill -HUP $MAINPID
ExecStop=/bin/kill -TERM $MAINPID
Restart=always
RestartSec=10
StartLimitInterval=60
StartLimitBurst=3

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/logging /var/log/centralized /tmp
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictNamespaces=true

# Resource limits - CRITICAL FOR MEMORY CONTROL
LimitNOFILE=2048
LimitNPROC=256
MemoryMax=256M
MemoryHigh=200M
CPUQuota=50%
TasksMax=128

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=logging-dashboard

# Health monitoring
WatchdogSec=30

[Install]
WantedBy=multi-user.target
